# contrib/caudater.init.conf
description "regular background program processing daemon"

start on runlevel [345]
stop on runlevel [!345]

respawn

exec /usr/bin/caudater /etc/yandex/caudater/config.yaml
